# MK-mid-ros2/mk_mid_bringup/systemd/connect-can.service
[Unit]
After=network.service time-sync.target
[Service]
Type=simple
ExecStart=/usr/sbin/connect-can
[Install]
WantedBy=multi-user.target
