# MK-mid-ros2/mk_mid_bringup/systemd/nav2.service
[Unit]
After=NetworkManager.service time-sync.target

[Service]
Type=simple
User=ros2
ExecStart=/usr/sbin/nav2

[Install]
WantedBy=multi-user.target
