# MK-mid-ros2/mk_mid_bringup/systemd/update_sub.service
[Unit]
After=NetworkManager.service time-sync.target

[Service]
Type=simple
User=ros2
ExecStart=/usr/sbin/update_sub

[Install]
WantedBy=multi-user.target
